# deploy/print-bridge/print-bridge.service
[Unit]
Description=Alavont Print Bridge
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/print-bridge
ExecStart=/usr/bin/node /opt/print-bridge/server.js
Restart=always
RestartSec=5
EnvironmentFile=/opt/print-bridge/.env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
